# etc/supervisor/conf.d/couchpotatoserver.conf
[program:couchpotatoserver]
autorestart = true
autostart = true
command = couchpotatoserver.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/couchpotatoserver.log
stopasgroup = true
