# etc/supervisor/conf.d/couchpotatoserver.conf
[program:couchpotatoserver]
autorestart = true
autostart = true
command = couchpotatoserver.sh
stopasgroup = true
