# etc/supervisor/conf.d/daily.conf
[program:daily]
autorestart = true
autostart = true
command = daily.sh
priority = 1000
stopasgroup = true
