# etc/supervisor/conf.d/logentries.conf
[program:logentries]
autorestart = true
autostart = true
command = logentries.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/logentries.log
stopasgroup = true
