# etc/supervisor/conf.d/nginx.conf
[program:nginx]
autorestart = true
autostart = true
command = nginx.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/nginx.log
stopasgroup = true
