# etc/supervisor/conf.d/nginx.conf
[program:nginx]
autorestart = true
autostart = true
command = nginx.sh
stopasgroup = true
