# etc/supervisor/conf.d/nzbget.conf
[program:nzbget]
autorestart = true
autostart = true
command = nzbget.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/nzbget.log
stopasgroup = true
