# etc/supervisor/conf.d/nzbget.conf
[program:nzbget]
autorestart = true
autostart = true
command = nzbget.sh
stopasgroup = true
