# etc/supervisor/conf.d/rclone.conf
[program:rclone]
autorestart = true
autostart = true
command = rclone.sh
stopasgroup = true
