# etc/supervisor/conf.d/rclone.conf
[program:rclone]
autorestart = true
autostart = true
command = rclone.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/rclone.log
stopasgroup = true
