# etc/supervisor/conf.d/sickrage.conf
[program:sickrage]
autorestart = true
autostart = true
command = sickrage.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/sickrage.log
stopasgroup = true
