# etc/supervisor/conf.d/sickrage.conf
[program:sickrage]
autorestart = true
autostart = true
command = sickrage.sh
stopasgroup = true
