# etc/supervisor/conf.d/transmission.conf
[program:transmission]
autorestart = true
autostart = true
command = transmission.sh
stopasgroup = true
