# etc/supervisor/conf.d/transmission.conf
[program:transmission]
autorestart = true
autostart = true
command = transmission.sh
redirect_stderr = true
stdout_logfile = /opt/var/log/transmission.log
stopasgroup = true
